# qr_code_scanner/qr_scanner.service
[Unit]
Description=Service for qr scanner
After=network.target

[Service]
User=root
WorkingDirectory=/home/pi/core_boxes_service/qr_code_scanner
ExecStart=python3 /home/pi/core_boxes_service/qr_code_scanner/main.py
Restart=always

[Install]
WantedBy=multi-user.target
